# docker/services.d/frontend/run
#!/usr/bin/execlineb -P
with-contenv
s6-setuidgid app
node /app/frontend/build/index.js